# components/minion/minion_server.conf
[program:minion_server]
command=/handler/minion_server.py
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

