# components/minion/worker.conf
[program:worker_{{image}}]
command=/handler/worker.py {{image}} {{args}}
process_name=worker_{{image}}_%(process_num)s
autostart=true
numprocs={{numprocs}}
