# dockerfiles/registry/registry.conf
[program:registry]
command={{DOCKER}} -H {{DOCKER_SOCKET}} run --rm -p 5000:5000 registry
autostart=true
