# dockerfiles/minion/minion_actor.conf
[program:minion_actor]
command=/handler/start_actor.py MinionServer
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

