# dockerfiles/minion/server.conf
[program:server]
command=/handler/server.py
autostart=true
