# dockerfiles/minion/worker.conf
[program:worker]
command=/handler/worker.py {{image}} {{args}}
process_name=worker_%(process_num)s
autostart=true
numprocs={{numprocs}}
