# ansible/roles/adama_server/templates/adama.ini.j2
[program:adama]
command=/home/adama/adama/bin/adama_server.py
autostart=true
autorestart=true
environment=HOME="/home/adama"
stopasgroup=true
killasgroup=true
stopsignal=INT
