# dockerfiles/example/example.conf
[program:example]
command=/example.py
autostart=true
