# dockerfiles/serfnode/programs/example.conf
[program:example]
command=/bin/cat
process_name=cat_%(process_num)s
numprocs={{numprocs}}
autostart=true

