# dockerfiles/zookeeper/zookeeper.conf
[program:zookeeper]
command={{DOCKER}} -H {{DOCKER_SOCKET}} run --rm -p 2181:2181 -p 2888:2888 -p 3888:3888 --name {{name}} --net host --volumes-from {{volume}} jplock/zookeeper
autostart=true
